# systemd/easytype.service
[Unit]
Description=EasyType voice dictation
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/easytype
Restart=on-failure

[Install]
WantedBy=default.target
